# supervisor.example.conf
; /etc/supervisor/conf.d/example.com.conf 
[program:example.com]
command = /srv/example.com/virtualenv/bin/gunicorn_start ; Command to start app
user = example.com ; User to run as
stdout_logfile = /srv/example.com/logs/gunicorn_supervisor.log ; Where to write log messages
redirect_stderr = true ; Save stderr in the same log

